# catkin_ws/devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/ubuntu/Peach_Bot/catkin_ws/src/ros/rosbash/rostcsh"
