# catkin_ws/devel_isolated/rosbuild/etc/catkin/profile.d/10.rosbuild.sh
# generated from rosbuild/env-hooks/10.rosbuild.sh.em

# env variables in develspace
export ROS_ROOT="/home/ubuntu/Peach_Bot/catkin_ws/src/ros/rosbuild"
